# Chapter4/hello-cluster.service
[Unit]

[Service]
ExecStart=/usr/bin/bash -c "while true; do echo 'Hello Cluster'; sleep 1; done"

